# deploy/launchd/com.supragnosis.daemon.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  supragnosis standalone MCP daemon (launchd LaunchAgent).
  Label is the canonical com.supragnosis.daemon so the CLI can drive it
  (supragnosis restart/stop -> launchctl kickstart -k / bootout).
  Install: see deploy/README.md. Paths are absolute for user ashon.lee - other users should adjust.
  loopback only (Principle 17) - MCP http://127.0.0.1:7373/mcp, viewer http://127.0.0.1:7374.
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.supragnosis.daemon</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/ashon.lee/.local/bin/supragnosis</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <key>SUPRAGNOSIS_HTTP_ADDR</key>
        <string>127.0.0.1:7373</string>
        <key>SUPRAGNOSIS_VIZ_ADDR</key>
        <string>127.0.0.1:7374</string>
        <key>SUPRAGNOSIS_STORE</key>
        <string>cozo</string>
        <key>SUPRAGNOSIS_DATA_DIR</key>
        <string>/Users/ashon.lee/.supragnosis/db</string>
        <key>SUPRAGNOSIS_WORKSPACE</key>
        <string>supragnosis</string>
        <key>SUPRAGNOSIS_HOST</key>
        <string>ashon-mac</string>
        <key>SUPRAGNOSIS_EMBED</key>
        <string>hashing</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/ashon.lee/.supragnosis/log/supragnosis.out.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/ashon.lee/.supragnosis/log/supragnosis.err.log</string>
</dict>
</plist>
